# vendor/github.com/docker/docker/contrib/init/openrc/docker.initd
#!/sbin/openrc-run
# Copyright 1999-2013 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

command="$***REMOVED***DOCKERD_BINARY:-/usr/bin/dockerd***REMOVED***"
pidfile="$***REMOVED***DOCKER_PIDFILE:-/run/$***REMOVED***RC_SVCNAME***REMOVED***.pid***REMOVED***"
command_args="-p \"$***REMOVED***pidfile***REMOVED***\" $***REMOVED***DOCKER_OPTS***REMOVED***"
DOCKER_LOGFILE="$***REMOVED***DOCKER_LOGFILE:-/var/log/$***REMOVED***RC_SVCNAME***REMOVED***.log***REMOVED***"
DOCKER_ERRFILE="$***REMOVED***DOCKER_ERRFILE:-$***REMOVED***DOCKER_LOGFILE***REMOVED******REMOVED***"
DOCKER_OUTFILE="$***REMOVED***DOCKER_OUTFILE:-$***REMOVED***DOCKER_LOGFILE***REMOVED******REMOVED***"
start_stop_daemon_args="--background \
	--stderr \"$***REMOVED***DOCKER_ERRFILE***REMOVED***\" --stdout \"$***REMOVED***DOCKER_OUTFILE***REMOVED***\""

start_pre() ***REMOVED***
	checkpath -f -m 0644 -o root:docker "$DOCKER_LOGFILE"

	ulimit -n 1048576

	# Having non-zero limits causes performance problems due to accounting overhead
	# in the kernel. We recommend using cgroups to do container-local accounting.
	ulimit -u unlimited

	return 0
***REMOVED***
